# supervisor/backend.conf
[program:backend]
command=python3 service.py
